# examples/soc-llm/soc-l1-triage.timer
[Unit]
Description=Run SOC L1 triage overnight (off-peak GPU time)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
